# os/kiosk/blackroad-kiosk.service
[Unit]
Description=BlackRoad Chromium Kiosk
After=graphical.target systemd-user-sessions.service

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/startx /opt/blackroad/os/kiosk/start-kiosk.sh -- :0
Restart=always
RestartSec=5
Description=BlackRoad Kiosk
After=graphical.target

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/startx /usr/bin/chromium --kiosk --no-first-run --fast --fast-start --window-size=1920,1080 http://localhost/
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
