# os/systemd/blackroad.target
[Unit]
Description=BlackRoad OS Target
Wants=blackroad-compose.service
After=blackroad-compose.service
