# systemd/blackroad-mirror.service
[Unit]
Description=BlackRoad mirror updater

[Service]
Type=oneshot
ExecStart=/bin/bash /srv/blackroad/health/mirror_sync.sh
User=root
