# systemd/blackroad-snapshot.service
[Unit]
Description=BlackRoad snapshot

[Service]
Type=oneshot
ExecStart=/srv/blackroad-backups/snapshot_now.sh
