# systemd/blackroad-validate.timer
[Unit]
Description=Run BlackRoad validator every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=blackroad-validate.service

[Install]
WantedBy=timers.target
