# systemd/blackroad-snapshot.timer
[Unit]
Description=Nightly snapshot

[Timer]
OnCalendar=*-*-* 03:33:00
Persistent=true

[Install]
WantedBy=timers.target
