# systemd/blackroad-identity.service
[Unit]
Description=BlackRoad identity ping

[Service]
Type=oneshot
ExecStart=/usr/local/bin/identity_ping.sh
