# systemd/blackroad-api.service
[Unit]
Description=BlackRoad.io API (Express)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/blackroad-api
EnvironmentFile=/srv/blackroad-api/.env
ExecStart=/usr/bin/node server_full.js
Restart=always
RestartSec=2
User=www-data
Group=www-data
StandardOutput=journal
StandardError=journal
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
