# systemd/lucidia-math-autoheal.service
[Unit]
Description=Lucidia Infinity Math auto-heal watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lucidia-math-autoheal.sh

[Install]
WantedBy=multi-user.target
