# systemd/lucidia-math-autoheal.timer
[Unit]
Description=Run Lucidia Infinity Math auto-heal every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=60
Unit=lucidia-math-autoheal.service

[Install]
WantedBy=timers.target
