# services/systemd/blackroad-api.service
[Unit]
Description=BlackRoad API Service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/blackroad-api
Restart=always
User=pi
WorkingDirectory=/home/pi

[Install]
WantedBy=multi-user.target
