# etc/cron.d/blackroad-backup
0 2 * * * root /usr/local/bin/blackroad-autoheal.sh backup
