# etc/cron.d/prism-snapshot
0 3 * * * root /usr/local/bin/prism-snapshot.sh
