# ops/systemd/seed/lucidia-vector.service
<!-- FILE: /ops/systemd/seed/lucidia-vector.service -->
[Unit]
Description=Lucidia Vector Service
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/lucidia/docker
ExecStart=/usr/bin/docker compose up -d vector
ExecStop=/usr/bin/docker compose stop vector
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
