# blackroad-agent/services/systemd/blackroad-agent.path
[Unit]
Description=Watch for BlackRoad Agent updates

[Path]
PathChanged=/opt/blackroad-agent
Unit=blackroad-agent.service

[Install]
WantedBy=multi-user.target
