# ops/backup/backup.service
# <!-- FILE: /ops/backup/backup.service -->
[Unit]
Description=Restic backup

[Service]
Type=oneshot
EnvironmentFile=/srv/blackroad/ops/backup/restic.env
ExecStart=/usr/local/bin/restic backup /srv/blackroad
