# lucidia/lucidia.target
[Unit]
Description=Lucidia Stack
Requires=lucidia.service
