# tools/lucidia-trt.service
[Unit]
Description=Lucidia (TensorRT-LLM, OpenAI-compatible)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment="PYTHONUNBUFFERED=1"
Environment="PORT=8000"
# Adjust paths if you installed elsewhere:
ExecStart=/home/ubuntu/.venvs/trtllm/bin/trtllm-serve /home/ubuntu/blackroad-prism-console/models/merged/lucidia-neox-1.4b --port ${PORT} --host 0.0.0.0
Restart=always
RestartSec=5
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/blackroad-prism-console
Description=Lucidia TensorRT-LLM Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/lucidia
ExecStart=/usr/bin/trtllm-serve /opt/lucidia/models/merged/lucidia-neox-1.4b --host 0.0.0.0 --port 8000
Restart=on-failure
Environment=MAX_TOKENS=512

[Install]
WantedBy=multi-user.target
