# pi-cortex-stack/pi/services/pi-holo-renderer.service
[Unit]
Description=Pi Cortex Holo Renderer
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-cortex-stack
ExecStart=/home/pi/agent-venv/bin/python /home/pi/pi-cortex-stack/pi/pi_holo_renderer.py
Restart=on-failure
EnvironmentFile=/home/pi/pi-cortex-stack/.env

[Install]
WantedBy=multi-user.target
