# etc/systemd/system/blackroad-autoheal.timer
[Unit]
Description=Run BlackRoad Auto-Heal every minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Unit=blackroad-autoheal.service
Description=Run BlackRoad Auto-Heal every 1 min

[Timer]
OnBootSec=1min
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
