# etc/systemd/system/br-cleanup-nightly.timer
[Unit]
Description=Run br-cleanup nightly

[Timer]
OnCalendar=03:10
Persistent=true

[Install]
WantedBy=timers.target

