# etc/systemd/system/br-cleanup-nightly.service
[Unit]
Description=BlackRoad nightly cleanup (dry-run + prune)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/br-cleanup.sh audit
ExecStart=/usr/local/sbin/br-cleanup.sh prune

