# etc/systemd/system/chem-agent.service
[Unit]
Description=BlackRoad Lucidia Chemistry Spooler (PySCF)
After=network.target

[Service]
Type=simple
User=blackroad
WorkingDirectory=/opt/blackroad/agents/pyscf
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 /opt/blackroad/agents/pyscf/pyscf_spooler.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
