# etc/systemd/system/symbol-gateway.service


[Unit]
Description=BlackRoad Symbol Gateway (Typst Codex)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/blackroad/symbol-gateway
ExecStart=/opt/blackroad/symbol-gateway/target/release/symbol-gateway serve --addr 127.0.0.1:8087
Restart=on-failure
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
